# docker/rootfs/etc/services.d/wishyoo-web/run
#!/usr/bin/with-contenv bash

cd /app

# start our node.js application
npm start;
